# supervisor/goaccess_supervisor.conf
[program:goaccess_cronjob]
command=/usr/local/bin/long.sh
autostart=true
autorestart=true