# supervisor/letsencrypt_supervisor.conf
[program:letsencrypt_entrypoint]
command=/opt/nginx-letsencrypt/entrypoint.sh
autostart=true
autorestart=false